# crontab
*   *   *   *   *     cd /config; shreddit >> /var/log/cron.log 2>&1
